# home/.chezmoitemplates/launchd-agent.plist.tmpl
{{- /* Input: dict with label (string), args (list), log (string), workdir (string) */ -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>{{ get . "label" | default "" }}</string>
  <key>KeepAlive</key>
  <false/>
  <key>RunAtLoad</key>
  <false/>
  <key>LimitLoadToSessionType</key>
  <array>
    <string>Aqua</string>
    <string>Background</string>
    <string>LoginWindow</string>
    <string>StandardIO</string>
    <string>System</string>
  </array>
  <key>ProgramArguments</key>
  <array>
{{- range (get . "args" | default list) }}
    <string>{{ . }}</string>
{{- end }}
  </array>
{{- with (get . "log") }}
  <key>StandardErrorPath</key>
  <string>{{ . }}</string>
  <key>StandardOutPath</key>
  <string>{{ . }}</string>
{{- end }}
{{- with (get . "workdir") }}
  <key>WorkingDirectory</key>
  <string>{{ . }}</string>
{{- end }}
</dict>
</plist>
