# home/Library/private_LaunchAgents/com.process-compose.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.process-compose</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{ .brewPrefix }}/bin/process-compose</string>
    <string>up</string>
    <string>--use-uds</string>
    <string>--tui=false</string>
    <string>--ordered-shutdown</string>
    <string>--logs-truncate</string>
    <string>--hide-disabled</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>{{ .chezmoi.homeDir }}/.local/bin:{{ .brewPrefix }}/bin:{{ .brewPrefix }}/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>{{ .chezmoi.homeDir }}</string>
    <key>PC_CONFIG_FILES</key>
    <string>{{ .chezmoi.homeDir }}/.config/process-compose/global.yaml</string>
    <key>PC_LOG_FILE</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/process-compose/process-compose.log</string>
    <key>PC_SOCKET_PATH</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/process-compose/global.sock</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/.local/state/process-compose/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/.local/state/process-compose/stderr.log</string>
</dict>
</plist>
